# server/scripts/airastro.service
[Unit]
Description=AirAstro Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/airastro/current/server
ExecStartPre=/opt/airastro/scripts/setup-current-version.sh
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
